# deploy/webaeger.service
[Unit]
Description=Webaeger (Grill Master) temperature server
After=network.target

[Service]
Type=simple
User=REPLACE_USER
WorkingDirectory=/opt/Webaeger
Environment=NODE_ENV=production
Environment=SENSOR=rtd
Environment=PORT=3000
Environment=POLL_MS=30000
Environment=TEMP_UNIT=F
Environment=DB_PATH=/opt/Webaeger/data/grillmaster.db
ExecStart=/usr/bin/npm run start --workspace=@grill-master/server
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
